# examples/dockerapp/supervisord-videoai.conf
[program:videoai]
command=python3 /app/my_video_detect_client.py
numprocs=1
autostart=true
autorestart=false
user=root