# examples/dockerapp/supervisord-videoai.conf
[program:videoai]
command=python3 /app/my_video_detect_client.py
numprocs=1
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/app/app.log
user=root